# always-on/qaio-engine.service
[Unit]
Description=Qaio Engine
After=network.target

[Service]
Type=simple
User=qaio
Group=qaio
WorkingDirectory=/opt/qaio
Environment=QAIO_BIND=0.0.0.0:7777
Environment=QAIO_BIND_ALL=1
EnvironmentFile=/etc/qaio/engine.env
ExecStart=/usr/local/bin/qaio-engine
Restart=on-failure
RestartSec=3s
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/var/lib/qaio

[Install]
WantedBy=multi-user.target
